# backend/outeniqua-contact.service
[Unit]
Description=Outeniqua Studios Contact Form API
After=network.target

[Service]
User=www-data
WorkingDirectory=/var/www/outeniqua-b2b/backend
ExecStart=/var/www/outeniqua-b2b/backend/venv/bin/gunicorn --workers 3 --bind 127.0.0.1:5000 wsgi:app
Restart=always

[Install]
WantedBy=multi-user.target
